# server/crontab-main/league-waivers.cron
# Waivers - Poaching

## Process Poaching Waivers Every 15 Minutes (except at 15:00)
*/15 0-14,16-23 * * * NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/process-poaching-waivers.mjs >> /var/log/league/league-process-poaching-waivers.log 2>&1

## Process Poaching Waivers Daily at 15:15
15 15 * * * NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/process-poaching-waivers.mjs --daily >> /var/log/league/league-process-poaching-waivers-daily.log 2>&1



# Waivers

## At 15:00 every day in January, February, and every month from August through December
0 15 * 1,2,8-12 * NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/process-waivers-free-agency.mjs --daily >> /var/log/league/league-waivers-free-agency.log 2>&1

## At every 5th minute past every hour from 15 through 23 on Wednesday.
5-59/5 15-23 * * 3 NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/process-waivers-free-agency.mjs >> /var/log/league/league-waivers-free-agency.log 2>&1

## At every 5th minute from 5 through 59 on Monday and every day-of-week from Thursday through Sunday.
5-59/5 * * * 1,4-7 NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/process-waivers-free-agency.mjs >> /var/log/league/league-waivers-free-agency.log 2>&1



# Poaching Claims

## Process Poaching Claims Every Hour
0 */1 * * * NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/process-poaching-claims.mjs >> /var/log/league/league-process-poaching-claims.log 2>&1
