# server/crontab-main/league-backups.cron
# Checkpoint Backup

## Backup User Data (At minute 0 in every month from June through August) (Checkpoint)
0 * * 6-8 * job-wrapper.sh league-backup-user-checkpoint /bin/bash /root/league/scripts/postgres-backup.sh -n checkpoint >> /var/log/league/league-postgres-backup.log 2>&1

## Backup User Data (At every 15th minute in Jan-Feb and Sep-Dec) (Checkpoint)
*/15 * * 1,2,9-12 * job-wrapper.sh league-backup-user-checkpoint /bin/bash /root/league/scripts/postgres-backup.sh -n checkpoint >> /var/log/league/league-postgres-backup.log 2>&1

## Backup Stat Data (Once a day) (Checkpoint)
0 2 * 9-12 3 job-wrapper.sh league-backup-stats-checkpoint /bin/bash /root/league/scripts/postgres-backup.sh -s -n checkpoint >> /var/log/league/league-postgres-backup-stats.log 2>&1

## Backup Betting Data (Once a day) (Checkpoint)
30 2 * 9-12 3 job-wrapper.sh league-backup-betting-checkpoint /bin/bash /root/league/scripts/postgres-backup.sh -b -n checkpoint >> /var/log/league/league-postgres-backup-betting.log 2>&1

## Backup Projections Data (Once a day) (Checkpoint)
0 3 * 9-12 3 job-wrapper.sh league-backup-projections-checkpoint /bin/bash /root/league/scripts/postgres-backup.sh -p -n checkpoint >> /var/log/league/league-postgres-backup-projections.log 2>&1

## Backup Log Data (Once a day) (Checkpoint)
30 3 * 9-12 3 job-wrapper.sh league-backup-log-checkpoint /bin/bash /root/league/scripts/postgres-backup.sh -l -n checkpoint >> /var/log/league/league-postgres-backup-log.log 2>&1

## Backup Cache Data (Once a day) (Checkpoint)
0 4 * 9-12 3 job-wrapper.sh league-backup-cache-checkpoint /bin/bash /root/league/scripts/postgres-backup.sh -c -n checkpoint >> /var/log/league/league-postgres-backup-cache.log 2>&1



# Time-Series Backup

## Backup User Data (Once a day) (Time-Series)
0 0 * * * job-wrapper.sh league-backup-user-timeseries /bin/bash /root/league/scripts/postgres-backup.sh >> /var/log/league/league-postgres-backup.log 2>&1

## Backup Stat Data (At 00:00 on Wednesday in every month from September through December) (Time-Series)
0 0 * 9-12 3 job-wrapper.sh league-backup-stats-timeseries /bin/bash /root/league/scripts/postgres-backup.sh -s >> /var/log/league/league-postgres-backup-stats.log 2>&1

## Backup All Betting Data (At 06:00 on Wednesday in every month from September through December) (Time-Series)
0 6 * 9-12 3 job-wrapper.sh league-backup-betting-timeseries /bin/bash /root/league/scripts/postgres-backup.sh -b >> /var/log/league/league-postgres-backup-betting.log 2>&1

## Backup All Projections Data (At 08:00 on Wednesday in every month from September through December) (Time-Series)
0 8 * 9-12 3 job-wrapper.sh league-backup-projections-timeseries /bin/bash /root/league/scripts/postgres-backup.sh -p >> /var/log/league/league-postgres-backup-projections.log 2>&1

## Backup All Cache Data (At 10:00 on Wednesday in every month from September through December) (Time-Series)
0 10 * 9-12 3 job-wrapper.sh league-backup-cache-timeseries /bin/bash /root/league/scripts/postgres-backup.sh -c >> /var/log/league/league-postgres-backup-cache.log 2>&1

## Backup All Log Data (At 12:00 on Wednesday in every month from September through December) (Time-Series)
0 12 * 9-12 3 job-wrapper.sh league-backup-log-timeseries /bin/bash /root/league/scripts/postgres-backup.sh -l >> /var/log/league/league-postgres-backup-log.log 2>&1
