# server/crontab-main/league-maintenance.cron
# League

## Every day at 4:00AM
0 4 * * * NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/calculate-team-daily-ktc-value.mjs >> /var/log/league/league-calculate-team-daily-ktc-value.log 2>&1

## Every day at 4:30AM
30 4 * * * NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/generate-league-team-player-seasonlogs.mjs --all >> /var/log/league/league-generate-league-team-player-seasonlogs.log 2>&1

## Every Monday at 5:00AM (with optimal-lineup pass)
0 5 * * 1 NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/generate-league-team-player-seasonlogs.mjs --all --with-optimal >> /var/log/league/league-generate-league-team-player-seasonlogs-with-optimal.log 2>&1

## Every day at 5:00AM (--rebuild reasserts canonical state; the walker is
## idempotent and lid=1 rebuild runs in ~5 seconds, so the daily cost is
## negligible. Without --rebuild any future emission bug that leaks past
## the suppression guards in walk-transactions.mjs accumulates orphans
## the partial unique indexes cannot recover from on their own.)
0 5 * * * NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/generate-roster-asset-lineage.mjs --all --rebuild >> /var/log/league/league-generate-roster-asset-lineage.log 2>&1

## Calibrate composite market value sources (15 5 daily; 15-min buffer after FantasyPros draft import at 0 5 * 5-8 *)
15 5 * * * NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/calibrate-composite-market-value-sources.mjs >> /var/log/league/calibrate-composite-market-value-sources.log 2>&1

## Generate composite market value daily (45 5; runs after calibration)
45 5 * * * NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/generate-composite-market-value.mjs >> /var/log/league/generate-composite-market-value.log 2>&1

## Announce restricted free agent at 9pm
0 21 * * * NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/announce-restricted-free-agent.mjs >> /var/log/league/announce-restricted-free-agent.log 2>&1

## Process restricted free agency bids at 6pm
0 18 * * * NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/process-restricted-free-agency-bids.mjs >> /var/log/league/process-restricted-free-agency-bids.log 2>&1

## Reset stale RFA tags (defense-in-depth) — daily at 5am Feb-May. Idempotent: preserves any tag=4 backed by an un-cancelled current-year bid in restricted_free_agency_bids.
0 5 * 2-5 * NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/reset-player-restricted-free-agency-tags.mjs >> /var/log/league/reset-player-restricted-free-agency-tags.log 2>&1

## Announce free agency period start (check every 30 minutes in August and September)
*/30 * * 8,9 * NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/announce-free-agency-period-start.mjs >> /var/log/league/announce-free-agency-period-start.log 2>&1

## Process extensions at ext_date (check every 5 minutes in May; self-gating fires once per league per ext_date via league_notifications marker)
*/5 * * 5 * NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/process-extensions.mjs >> /var/log/league/process-extensions.log 2>&1

# League Each Season Week

## At 23:55 on Monday during the NFL Season
55 23 * 1,2,9-12 1 NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/generate-rosters.mjs >> /var/log/league/league-generate-rosters.log 2>&1



# Maintenance

## Generate new season records (At 07:00 on the 1st, Feb through Aug - no-op until season-dates.mjs is updated)
0 7 1 2-8 * NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/generate-seasons.mjs >> /var/log/league/generate-seasons.log 2>&1

## Create database partitions for upcoming season (At 00:00 on the 1st in May and July)
0 0 1 5,7 * NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/create-season-partitions.mjs >> /var/log/league/create-season-partitions.log 2>&1

## Every day at 0:00
0 0 * * * job-wrapper.sh league-cleanup-puppeteer find /tmp -maxdepth 1 -name 'puppeteer*' -mtime +2 -exec rm -rf {} \; >> /var/log/league/league-cleanup-puppeteer.log 2>&1

## Check disk space every hour
0 * * * * NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/check-disk-space-main-server.mjs >> /var/log/league/league-check-disk-space-main-server.log 2>&1



# Crontab Export

## Export crontab for backup sync (hourly)
0 * * * * crontab -l > /tmp/league-crontab.txt
