# server/crontab-worker-1.cron
################################### GAMES ########################################

## Update NFL Schedule At 03:00 every day (Regular Season)
0 3 * 1,2,9-12 * NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/import-nfl-games-ngs.mjs &>> /var/log/league/import-nfl-games-ngs.log
0 3 * 1,2,9-12 * NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/import-nfl-games-nfl.mjs --current --force &>> /var/log/league/league-import-nfl-games-nfl.log

## Update NFL Schedule At 03:00 every day (Off Season)
0 3 * 5-8 0 NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/import-nfl-games-ngs.mjs &>> /var/log/league/import-nfl-games-ngs.log

## At every 15th minute past every hour from 19 through 20 on Monday in January, February, and every month from September through December.
# */15 19-20 * 1,2,9-12 1 NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/import-nfl-games-nfl.mjs --current --force &>> /var/log/league/league-import-nfl-games-nfl.log

## At every 15th minute past every hour from 19 through 20 on Thursday in January, February, and every month from September through December.
# */15 19-20 * 1,2,9-12 4 NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/import-nfl-games-nfl.mjs --current --force &>> /var/log/league/league-import-nfl-games-nfl.log

## At every 15th minute on Sunday in January, February, and every month from September through December.
# */15 * * 1,2,9-12 7 NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/import-nfl-games-nfl.mjs --current --force &>> /var/log/league/league-import-nfl-games-nfl.log

# DFS Salaries

## Every Wednesday at 1:20PM
20 13 * * 3 NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/import-draftkings-salaries.mjs &>> /var/log/league/league-import-draftkings-salaries.log

################################### WEEK FINALIZATION ########################################

## Finalize previous week at midnight Tuesday (primary)
0 0 * 1,2,9-12 2 NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/jobs/finalize-week.mjs &>> /var/log/league/finalize-week.log

## Finalize previous week at midnight Wednesday (backup)
0 0 * 1,2,9-12 3 NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/jobs/finalize-week.mjs &>> /var/log/league/finalize-week.log

################################### SEASON FINALIZATION ########################################

## Re-finalize season at 6AM on January 15th (approximately 1 week after season end, for stat corrections)
0 6 15 1 * NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/finalize-season.mjs &>> /var/log/league/finalize-season.log

## Re-finalize season at 6AM on February 15th (approximately 1 month after season end, for final stat corrections)
0 6 15 2 * NODE_ENV=production /root/.nvm/versions/node/current/bin/node /root/league/scripts/finalize-season.mjs &>> /var/log/league/finalize-season.log
